# config/unitmail-gateway.service
[Unit]
Description=unitMail Gateway Service
Documentation=https://unitmail.io/docs
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=exec
User=unitmail
Group=unitmail

# Working directory
WorkingDirectory=/opt/unitmail

# Environment file support
EnvironmentFile=-/etc/unitmail/gateway.env
EnvironmentFile=-/opt/unitmail/.env

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="UNITMAIL_CONFIG_FILE=/etc/unitmail/gateway.toml"

# Main process
ExecStart=/opt/unitmail/venv/bin/python -m gateway.app
ExecReload=/bin/kill -HUP $MAINPID

# Startup and shutdown
TimeoutStartSec=30
TimeoutStopSec=30
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
LockPersonality=yes

# Allow binding to privileged ports (25, 465, 587)
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Read-write paths
ReadWritePaths=/var/lib/unitmail
ReadWritePaths=/var/log/unitmail
ReadWritePaths=/var/spool/unitmail

# System paths
ReadOnlyPaths=/etc/unitmail

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=unitmail-gateway

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Process settings
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
